# systemd/claudius-api.service
[Unit]
Description=Claudius API Server
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=claudius
Group=claudius
WorkingDirectory=/opt/claudius
ExecStart=/usr/bin/python3 /opt/claudius/claudius-api-with-telegram.py --port 3100

# IMPROVED: Always restart (not just on-failure)
Restart=always
RestartSec=5

# IMPROVED: More restart attempts before giving up
StartLimitIntervalSec=300
StartLimitBurst=10

# Resource limits - generous allocation for thinking mode + MCP servers
MemoryHigh=6G
MemoryMax=8G
CPUQuota=400%
LimitNOFILE=65536
TasksMax=500

# Priority boost - Claudius gets preferential treatment
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Watchdog disabled - Claudius doesn't implement sd_notify heartbeats
# WatchdogSec=120
# NotifyAccess=all

# Environment
Environment=PATH=/usr/local/bin:/usr/bin:/bin
EnvironmentFile=-/opt/omniops/.env

# Logging
StandardOutput=append:/var/log/claudius-api.log
StandardError=append:/var/log/claudius-api.log

# ADDED: Kill child processes on stop
KillMode=control-group
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
